# deploy/systemd/yanghub.service
[Unit]
Description=YangHub Django API
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/Yang_Hub/backend
EnvironmentFile=/var/www/Yang_Hub/.env
ExecStart=/var/www/Yang_Hub/.venv/bin/gunicorn backend.wsgi:application --bind 127.0.0.1:8000 --workers 2 --timeout 120
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
